# .gemini/antigravity/brain/c718aa18-28de-4cd5-9403-f4be6c1ae8db/Job_Search_Agent_System/crontab
# Job Search Agent — crontab pour déploiement manuel sur Contabo
# 1) Toutes les 2h : scan + pipeline full (mode both)
# 2) Quotidien 8h : relances J+4/J+10
#
# Installation : crontab crontab
#
0 */2 * * * cd /app && python -m scheduler.cron_runner --mode both --sources wttj,francetravail,indeed,hellowork,dogfinance,meteojob,glassdoor,linkedin,apec --max 10 >> /app/logs/cron.log 2>&1
0 8 * * * cd /app && python -m scheduler.followup_runner >> /app/logs/followup.log 2>&1
